# nullidentdmod.socket
[Unit]
Description=Socket for identd (NullidentdMod)

[Socket]
ListenStream=113
Accept=yes

[Install]
WantedBy=sockets.target
